# systemd/scraper-reload.timer
[Unit]
Description=TeslaScraper Reload timer

[Timer]
OnBootSec=50m
OnUnitActiveSec=50m

[Install]
WantedBy=timers.target