# systemd/scraper-archive.timer
[Unit]
Description=TeslaScraper Archive timer

[Timer]
OnBootSec=15s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target